# systemd/aniworld-sync.service
[Unit]
Description=AniWorld Sync Service
After=network.target aniworld-api.service aniworld-metadata.service

[Service]
Type=oneshot
User=emby
Group=emby
ExecStart=/opt/aniworld/venv/bin/python3 /opt/aniworld/sync.py
Environment=ANIWORLD_CONFIG=/etc/aniworld/config.ini
